# etc/rc.local
#!/bin/sh -e
rfkill block bluetooth
# echo PS2M > /proc/acpi/wakeup
exit 0
