# etc/systemd/system/ip4tables.service
[Unit]
Description=Automatically load /etc/ip4tables.rules
Before=network-pre.target

[Service]
Type=oneshot
ExecStart=/sbin/iptables-restore /etc/ip4tables.rules
RemainAfterExit=no
TimeoutSec=0

[Install]
RequiredBy=networking.service
